# infra/systemd/ytd-worker.service
[Unit]
Description=YT Downloader Worker
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=local-lenovo
WorkingDirectory=/home/local-lenovo/yt-downloader
EnvironmentFile=/home/local-lenovo/yt-downloader/apps/backend/.env
ExecStart=/home/local-lenovo/yt-downloader/bin/ytd-worker
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
